# lib/systemd/podman-dump@.timer
[Unit]
Description=Container backups %i schedule

[Timer]
OnCalendar=%i
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
